# ServiceChecker/path/service_verify_path.path
[Unit]
Description=Detect systemd Path

[Path]
PathChanged=/opt/tem_monitor
MakeDirectory=yes
Unit=service_verify_path.service

[Install]
WantedBy=multi-user.target